# n8n/rootfs/etc/services.d/n8n/run
#!/usr/bin/with-contenv bashio

# Expose n8n locally
export N8N_PORT=5678
export N8N_HOST=0.0.0.0

# Ensure data persistence
export N8N_USER_FOLDER=/data

# Setup integration with Home Assistant
export GENERIC_TIMEZONE=$(bashio::info.timezone)
export N8N_EDITOR_BASE_URL=$(bashio::addon.ingress_url)
export N8N_LOG_LEVEL=$(bashio::config 'log_level')
export N8N_PATH="$(bashio::addon.ingress_entry | sed 's|/*$|/|')"
export N8N_RUNNERS_ENABLED=$(bashio::config 'enable_task_runners')
export N8N_SECURE_COOKIE=$(bashio::config 'enable_ssl')

EXT_MODULES=""
for package in $(bashio::config 'node_function_external_modules'); do
    EXT_MODULES="${EXT_MODULES}${package},"
done
export NODE_FUNCTION_ALLOW_EXTERNAL=${EXT_MODULES%,}

# Start n8n
cd /usr/src/n8n
export PATH="/root/.local/share/mise/shims:$PATH"
exec ./node_modules/.bin/n8n start
