# n8n/rootfs/etc/services.d/n8n/run
#!/usr/bin/with-contenv bashio

exec /usr/bin/n8n
